# legacy/apache-php/bin/php74-mssql-ldap/rootfs/etc/services.d/apache2/run
#!/usr/bin/with-contenv sh

exec docker-php-entrypoint apache2-foreground

